# disable_lid_button.service
[Unit]
Description=Disable lid button
 
[Service]
ExecStart=/opt/ar18/disable_lid_button/disable_lid_button.sh
 
[Install]
WantedBy=multi-user.target